# cron/crontab
PATH='/usr/local/bin/:/usr/bin:/bin'
MAILTO='andre.lersveen@nrk.no'

00 8-17 * * * ~/.ssh/update_nrk_config.sh
